# deploy/systemd/rfp-frontend.service
# RFP RAG Streamlit frontend.
# The only externally exposed service (GCP firewall: tcp:8501).
# Server config (port 8501, address 0.0.0.0, 20MB upload) comes from
# frontend/.streamlit/config.toml, hence WorkingDirectory=frontend/.
#
# 설치 (경로·계정이 다르면 아래 User/WorkingDirectory/ExecStart 수정):
#   sudo cp deploy/systemd/rfp-frontend.service /etc/systemd/system/
#   sudo systemctl daemon-reload && sudo systemctl enable --now rfp-frontend

[Unit]
Description=RFP RAG Frontend (Streamlit, 0.0.0.0:8501)
After=network.target rfp-api.service
Wants=rfp-api.service

[Service]
Type=simple
User=seen
WorkingDirectory=/home/seen/intermediate-project/frontend
Environment=RAG_API_BASE_URL=http://127.0.0.1:8080
ExecStart=/home/seen/venvs/rfp/bin/streamlit run app.py
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
